# systemd/mem0-sync.service
[Unit]
Description=Sync voice agent conversations to Mem0
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=kevin
WorkingDirectory=/home/kevin/coding/voice-agent
ExecStart=/home/kevin/.local/bin/uv run python scripts/nightly-mem0-sync.py
StandardOutput=journal
StandardError=journal
